# docker/supervisord/vite.conf
[program:vite]
command=npm run dev -- --host
directory=/var/www/app
autostart=true
autorestart=true
stderr_logfile=/var/log/vite.err.log
stdout_logfile=/var/log/vite.out.log
user=www-data
